# ansible/roles/sgoinfre_mount/templates/systemd/sgoinfre.mount
[Unit]
After=network-online.target

[Mount]
What=sgoinfre:/export/sgoinfre
Where=/sgoinfre
Options=nfsvers=4,nolock,noatime

[Install]
WantedBy=multi-user.target
