# ansible/roles/udbsync_rfs/templates/systemd/rootssh.path
[Unit]
Description=Copy ~root/.ssh/authorized_keys to NFS root (watcher)

[Path]
PathChanged=/root/.ssh/authorized_keys
Unit=rootssh-copy.service

[Install]
WantedBy=multi-user.target
