# etc/systemd/system/teeworlds-server.service
[Unit]
Description=Teeworlds server

[Service]
DynamicUser=yes
ExecStart=/usr/bin/teeworlds_srv -f /etc/teeworlds_srv.cfg
